# .config/autostart/google-drive-ocamlfuse.desktop
[Desktop Entry]
Exec=google-drive-ocamlfuse ~/migoogledrive
Icon=
Name=google-drive-ocamlfuse
Path=
Terminal=False
Type=Application
