# .config/autostart/guake.desktop
[Desktop Entry]
Exec=guake
Icon=
Name=Guake
Path=
Terminal=False
Type=Application
