# .config/autostart/muezzin.desktop
[Desktop Entry]
Exec=muezzin
Icon=
Name=Muezzin
Path=
Terminal=False
Type=Application
