# .config/autostart/psensor.desktop
[Desktop Entry]
Exec=psensor
Icon=
Name=PSensor
Path=
Terminal=False
Type=Application
