# .config/autostart/StartNODE.desktop
[Desktop Entry]
Exec=node /home/mustafaozver/startkit/start.js
Icon=
Name=StartNODE
Path=
Terminal=False
Type=Application
