# .config/autostart/startsh.desktop
[Desktop Entry]
Exec=sh /home/mustafaozver/startsh.sh
Icon=
Name=StartSH
Path=
Terminal=False
Type=Application
